# build/openshift-origin-broker/init-scripts/livesys-late-openshift
#!/bin/bash
#
# live: Late init script for live image
#
# chkconfig: 345 99 01
# description: Late init script for configuring image

start() {
  IP_DEVICE=${IP_DEVICE:=eth0}
  IPADDR=$(ip addr show dev ${IP_DEVICE} | awk '/inet / { split($2,a, "/") ; print a[1];}')
  /usr/bin/ss-setup-node -b $IPADDR -n broker
  
  echo "initialize mongodb"
  service mongod start
  
  #wait for "[initandlisten] waiting for connections"
  WAIT=1
  while [ 1 -eq $WAIT ] ; do /bin/fgrep "[initandlisten] waiting for connections" /var/log/mongodb/mongodb.log && WAIT=$? ; echo "Waiting for mongo to initialize ...\n" ; tail -5 /var/log/mongodb/mongodb.log ; sleep 1 ; done
  
  /usr/bin/mongo localhost/stickshift_broker_dev --eval 'db.addUser("stickshift", "mooo")'
  
  cat <<EOF >/etc/resolv.conf
  search example.com
  nameserver 127.0.0.1
EOF
  
  # Set up the initial forwarder
  echo "forwarders { 8.8.8.8 ; 8.8.4.4 ; } ;" > /var/named/forwarders.conf
  
  # set SELinux label for forwarders file
  /sbin/restorecon -v /var/named/forwarders.conf
  
  grep -l NM_CONTROLLED /etc/sysconfig/network-scripts/ifcfg-* | xargs perl -p -i -e '/NM_CONTROLLED/ && s/yes/no/i'
  
  #Register admin user
  mongo stickshift_broker_dev --eval 'db.auth_user.update({"_id":"admin"}, {"_id":"admin","user":"admin","password":"2a8462d93a13e51387a5e607cbd1139f"}, true)'
  
  service named restart
  /usr/bin/ss-register-dns -h broker -n $IPADDR
  
  service httpd restart
  service sshd restart
  service qpidd restart
  service mcollective restart
  lokkit --service=ssh
  lokkit --service=http
  lokkit --service=https
  lokkit --service=dns
  lokkit -p 5672:tcp
}

case "$1" in
  start)
	start
	;;
  stop)
	;;
esac

exit 0
